# repo_content/Source/gigagrok-bot/gigagrok.service
[Unit]
Description=GigaGrok Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gigagrok
Group=gigagrok
WorkingDirectory=/opt/gigagrok
ExecStart=/opt/gigagrok/venv/bin/python main.py
Restart=always
RestartSec=10
TimeoutStopSec=30
Environment=PYTHONUNBUFFERED=1

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/gigagrok /opt/gigagrok/workspaces
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gigagrok

[Install]
WantedBy=multi-user.target
